# source/usr/lib/systemd/system/netcfg-auto-wireless-ng-suspend.service
[Unit]
Description=Suspend netctl-auto-wireless-ng upon sleep
PartOf=sleep.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/netctl-auto-wireless suspend
ExecStop=/usr/bin/netctl-auto-wireless resume
KillMode=none
